# rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/bin/sh
mkdir -p /run/dbus/ /var/run/dbus/
dbus-daemon --system --nofork --nosyslog
